# scripts/etl/sejung-etl.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>세정학원 CRM · Aca2000 MSSQL → Supabase ETL 매 시간 동기화 (KST)</Description>
    <URI>\Sejung\CRM-ETL-Hourly</URI>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-05-22T00:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <Repetition>
        <Interval>PT1H</Interval>
        <Duration>P1D</Duration>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <!-- 이전 실행이 안 끝났으면 새 실행 무시 (중복 발사 차단) -->
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>false</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>true</WakeToRun>
    <!-- ETL 1회 평균 5~15분. 다음 시간 전에 무조건 끝나도록 55분 한도. -->
    <ExecutionTimeLimit>PT55M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>cmd.exe</Command>
      <!-- README-windows.md 안내대로 C:\sejung-crm 에 clone 했다면 import 만 하면 됨.
           다른 경로면 아래 두 줄의 C:\sejung-crm 를 실제 경로로 바꿀 것. -->
      <Arguments>/c "C:\sejung-crm\scripts\etl\run_all.bat"</Arguments>
      <WorkingDirectory>C:\sejung-crm</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
